# provisioning/roles/weather-hat/files/weather-hat.conf
description "Weather-hat"
author "daniel.fox@ukgateway.net"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec python /src/foxy1982/RaspberryPi/weather-hat/main.py
